# docker/images/cron/configs/crontab
# Default Laravel shedule
*/1 * * * * php /var/www/html/app/scripts/email.php >> /var/www/html/app/scripts/cron.log 2>&1
